# supervisor/worker.conf
[program:tradeworker]
chdir=/app
autoreload=true
command=celery -A sphinxquant worker --loglevel=INFO --concurrency=10 -n trade@%h